# dist/hercules-server.service
[Unit]
Description=Hercules Server
BindsTo=hercules-monitor.service
After=hercules-monitor.service

[Service]
ExecStart=/home/marco/hercules-server
StandardOutput=syslog
StandardError=syslog
WorkingDirectory=/home/marco
TimeoutSec=5
