# dist/hercules-server.service
[Unit]
Description=Hercules Server
BindsTo=hercules-monitor.service
After=hercules-monitor.service

[Service]
ExecStart=/usr/local/bin/hercules-server
StandardOutput=syslog
StandardError=syslog
WorkingDirectory=/tmp/
TimeoutSec=5
LimitCORE=infinity
