# dist/hercules-monitor.service
[Unit]
Description=Hercules Monitor
BindsTo=hercules-server.service

[Service]
ExecStart=/usr/local/bin/hercules-monitor
StandardOutput=syslog
StandardError=syslog
WorkingDirectory=/tmp/
TimeoutSec=5
